# elec2020/form/supervisor/conf.d/web.ini
[program:web]
process_name=%(program_name)s_%(process_num)02d
command=php artisan swoole:http start
autostart=true
numprocs=%(ENV_SUPERVISOR_WEB_COUNT)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
